# systemd/pi-temp-alerter.service
[Unit]
Description=Pi Temperature Alerter - System health monitoring
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi-temp-alerter
Group=pi-temp-alerter
WorkingDirectory=/opt/pi-temp-alerter
ExecStart=/opt/pi-temp-alerter/venv/bin/python main.py start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/pi-temp-alerter/logs /opt/pi-temp-alerter/data

[Install]
WantedBy=multi-user.target
