# systemd/pimon.service
[Unit]
Description=PiMon - System health monitoring
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=pimon
Group=pimon
WorkingDirectory=/opt/pimon
ExecStart=/opt/pimon/venv/bin/python main.py start
Restart=on-failure
RestartSec=10
WatchdogSec=90
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/pimon/logs /opt/pimon/data

[Install]
WantedBy=multi-user.target
